# ops/example/systemd/api/api.socket
[Socket]
ListenStream = 80
#ListenStream = 443
BindIPv6Only = both
Service      = api.service

[Install]
WantedBy = sockets.target